# config/webgme.conf
# WebGME upstart configuration

description "Upstart script to launch WebGME forever"

# WebGME depends on MongoDB
start on started mongod
stop on stopping mongod

# attempt to restart the service when it stops
respawn
respawn limit 5 10

# use webgme as the working directory
chdir /home/vagrant/cpswt/cpswt-meta/

script
    exec /home/vagrant/.nvm/versions/node/v4.5.0/bin/node app.js
end script

